# deploy/alfred-dashboard.container
# Alfred Dashboard Quadlet container unit
#
# Deploy to mower:
#   scp alfred-dashboard.container <mower>:/etc/containers/systemd/alfred-dashboard.container
#   ssh <mower> "systemctl daemon-reload && systemctl start alfred-dashboard"

[Unit]
Description=Alfred Dashboard (status & control UI)
After=network-online.target sunray.service
Wants=network-online.target

[Container]
Image=ghcr.io/autoditac/alfred-dashboard:latest
ContainerName=alfred-dashboard
Network=host
Environment=SUNRAY_HOST=127.0.0.1
Environment=SUNRAY_PORT=80
Environment=PORT=3000
Environment=TZ=Europe/Copenhagen
AutoUpdate=registry
PodmanArgs=--init

[Service]
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target default.target
